# docker/rootfs/etc/services.d/swift-container/run
#!/bin/sh
source /etc/profile

# swift-container-server /etc/swift/container-server.conf
exec s6-setuidgid swift swift-init container restart --no-daemon
