# tools/playbooks/multinode_setup/templates/rc.local.j2
#!/bin/bash

mkdir -p /var/cache/swift
chown {{ ansible_user_id }}:{{ ansible_user_gid }} /var/cache/swift*
mkdir -p /var/run/swift
chown {{ ansible_user_id }}:{{ ansible_user_gid }} /var/run/swift

exit 0
