# docker/rootfs/etc/services.d/swift-proxy/run
#!/bin/sh
source /etc/profile

# swift-proxy-server /etc/swift/proxy-server.conf
exec s6-setuidgid swift swift-init proxy restart --no-daemon
