# docker/rootfs/etc/services.d/swift-object/run
#!/bin/sh
source /etc/profile

exec s6-setuidgid swift swift-init object restart --no-daemon
