# debian/swift-object.swift-object-replicator
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swift-object-replicator
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Swift object replicator server
# Description:       Object replicator server for swift.
### END INIT INFO

/usr/bin/swift-init object-replicator $1

