# debian/swift-container.swift-container-updater
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swift-container-updater
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Swift container updater server
# Description:       Container updater server for swift.
### END INIT INFO

/usr/bin/swift-init container-updater $1

