# debian/swift-account.swift-account-auditor
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swift-account-auditor
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Swift account auditor server
# Description:       Account auditor server for swift.
### END INIT INFO

/usr/bin/swift-init account-auditor $1

