# debian/swift-account.swift-account-reaper
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swift-account-reaper
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Swift account reaper server
# Description:       Account reaper for swift.
### END INIT INFO

/usr/bin/swift-init account-reaper $1

